# examples/cron.d/gmjobs-metrics-export
PATH=/sbin:/bin:/usr/sbin:/usr/bin
55 04 * * * root /opt/gmjobs-metrics-exporter/gmjobs-metrics-exporter.sh


